# etc/xdg/autostart/xclip.desktop
[Desktop Entry]
Name=XClip
Comment=Clipboard manager
Exec=xclip
Terminal=false
Type=Application
NoDisplay=true
